# package/python-deb/lib/systemd/system/python-deb.service
[Unit]
Description=python-deb

[Service]
ExecStart=/usr/bin/project_deb_bin

[Install]
WantedBy=multi-user.target
