# package/python-deb/DEBIAN/postinst
#!/bin/sh

systemctl enable python-deb
systemctl start python-deb
